# supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command =/usr/local/bin/uwsgi --http 0.0.0.0:8000 --wsgi-file  /usr/src/app/hello.py --callable app
